# themes/cambium/root/etc/uci-defaults/luci-theme-cambium
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Cambium=/luci-static/cambium
    set luci.main.mediaurlbase=/luci-static/cambium
	commit luci
EOF
exit 0
